# deploy/com.switchboard.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Switchboard control-plane API (FastAPI/uvicorn), loopback-only on port 8400
       (8000 is taken by another local project). Remote access goes through
       Tailscale (see deploy/README) — never bind this to 0.0.0.0 directly.
       __PROJECT_DIR__ is rewritten by install.sh.
       The caffeinate wrapper is LOAD-BEARING beyond sleep: a launchd agent's TCC
       "responsible process" is its executable, and a bare (non-Apple) python is
       denied ~/Documents — the worker proved the caffeinate-wrapped form reads the
       project dir fine while the unwrapped API crash-looped on pyvenv.cfg. -->
  <key>Label</key>
  <string>com.switchboard.api</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-i</string>
    <string>__PROJECT_DIR__/.venv/bin/uvicorn</string>
    <string>api.app:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8400</string>
    <string>--log-level</string>
    <string>warning</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/api.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/api.log</string>
</dict>
</plist>
